# config/systemd/delayed_job.service.erb
[Unit]
Description = Delayed Job

[Service]
Environment = PATH=/shims:/usr/local/bin:/usr/bin:/bin
Environment = RAILS_ENV=<%= fetch(:rails_env) %>
Environment = PWD=<%= current_path %>
Environment = QUEUES=default,paperclip,mailers

WorkingDirectory = <%= current_path %>

ExecStart  = /home/tumlino_production/.rbenv/shims/bundle exec bin/delayed_job -n 1 start
ExecStop   = /home/tumlino_production/.rbenv/shims/bundle exec bin/delayed_job -n 1 stop
ExecReload = /bin/kill -HUP $MAINPID

PIDFile = <%= shared_path %>/tmp/pids/delayed_job.pid
TimeoutStopSec = 22
Restart = always

[Install]
WantedBy = multi-user.target
